# config/systemd/iem-web-services-podman.service
[Unit]
Description=IEM Web Services (Podman)
Wants=network-online.target
After=network-online.target

[Service]
Environment="IMAGE=ghcr.io/akrherz/iem-web-services:latest"
Environment="CONTAINER_NAME=iem-web-services"
Environment="ENV_FILE=/etc/iem-web-services/iemws.env"
Environment="PGPASS_SECRET=iemws-pgpass"
Environment="PGPASS_UID=57439"
Environment="PGPASS_GID=57439"
Environment="PUBLISH_HOST=127.0.0.1"
Environment="PUBLISH_PORT=8000"
Environment="WORKERS=6"
Environment="MAX_REQUESTS=2000"
EnvironmentFile=-/etc/iem-web-services/iemws.env

ExecStartPre=-/usr/bin/podman rm -f ${CONTAINER_NAME}
ExecStartPre=/usr/bin/podman pull ${IMAGE}
ExecStartPre=/usr/bin/podman secret inspect ${PGPASS_SECRET}
ExecStart=/usr/bin/podman run --name ${CONTAINER_NAME} --replace --publish ${PUBLISH_HOST}:${PUBLISH_PORT}:8000 --volume /mesonet/data:/mesonet/data:ro --volume /mnt/mesonet2:/mnt/mesonet2:ro --volume /opt/bufkit:/opt/bufkit:ro --env-file ${ENV_FILE} --env PGPASSFILE=/run/secrets/pgpass --secret ${PGPASS_SECRET},type=mount,target=pgpass,uid=${PGPASS_UID},gid=${PGPASS_GID},mode=0600 --log-driver=journald ${IMAGE} micromamba run -n iemws gunicorn -w ${WORKERS} --graceful-timeout 60 -k iemws.worker.RestartableUvicornWorker -b 0.0.0.0:8000 --max-requests ${MAX_REQUESTS} --max-requests-jitter 500 --log-level warning iemws.main:app
ExecStop=/usr/bin/podman stop -t 30 ${CONTAINER_NAME}
ExecStopPost=-/usr/bin/podman rm -f ${CONTAINER_NAME}

Restart=always
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
